# webapp/simple.conf
[program:myserver]
command=gunicorn wikihub_worker:app -w 8 -b 0.0.0.0:80

[supervisord]
logfile=/nome/ubuntu/supervisord.log
loglevel=debug
user=root